# helpmein/frontend/supervisord.d/frontend.conf
[program:frontend-nuxt]
process_name=%(program_name)s_%(process_num)02d
command=npm start --prefix=/app
autostart=true
autorestart=true
redirect_stderr=true
